# observer/unit/input-method-observer.socket
[Unit]
Description=Input Method observer socket

[Socket]
ListenStream=6631
Accept=false

[Install]
WantedBy=sockets.target
